# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid rowan
chdir /home/rowan/sites/SITENAME/source

exec ../virtualenv/bin/gunicorn \
    --bind 0.0.0.1:5001 \
    run:app
